# systemd/bgaction.service
[Unit]
Description=Execute bgaction script

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/bgaction

